# bambuddy/rootfs/etc/services.d/nginx/finish
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Stop NGINX
# ==============================================================================

bashio::log.info 'NGINX proxy stopped.'
