# bambuddy/rootfs/etc/services.d/nginx/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start NGINX
# ==============================================================================

bashio::log.info 'Start NGINX proxy...'
exec nginx -g "daemon off;"
